# Library/LaunchAgents/net.makerbox.standup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>net.makerbox.standup</string>
        <key>ProgramArguments</key>
        <array>
            <string>/opt/mctk-standup/bin/standup.sh</string>
            <string>ADD_TOKEN_HERE</string>
        </array>
        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string>/sbin:/bin:/usr/sbin:/usr/bin:/opt/homebrew/bin:/opt/homebrew/sbin</string>
        </dict>
        <key>StartCalendarInterval</key>
        <dict>
            <key>Minute</key>
            <integer>55</integer>
        </dict>
        <key>StandardOutPath</key>
        <string>/opt/mctk-standup/log/stdout.log</string>
        <key>StandardErrorPath</key>
        <string>/opt/mctk-standup/log/stderr.log</string>
    </dict>
</plist>